# OS_Overlay/LastOS/Scripts/lastos.service
[Service]
ExecStart=/LastOS/Scripts/sudo_firstlogon_service.sh
[Install]
WantedBy=default.target
