# OS_Overlay/etc/rc.local
#!/bin/sh -e
apt update
rm /etc/rc.local
exit 0
